# lab-06-bash/systemd/site-check.service
[Unit]
Description=Перевірка доступності веб-сайту

[Service]
User=sc
ExecStart=/opt/scripts/check_site.sh
